# init/ntrc.service
[Unit]
Description=NITOS proxies
After=syslog.target

[Service]
Type=simple

[Install]
WantedBy=multi-user.target
